# systemd/drukmix.service
[Unit]
Description=DrukMix Moonraker Agent
After=network-online.target moonraker.service
Wants=network-online.target
Requires=moonraker.service

[Service]
Type=simple
User=drukos
Group=drukos
WorkingDirectory=/opt/drukmix
Environment=DRUKMIX_CONFIG=/etc/drukmix/drukmix.cfg
ExecStart=/opt/drukmix/.venv/bin/python /opt/drukmix/drukmix_agent.py
Restart=always
RestartSec=2
Nice=5

[Install]
WantedBy=multi-user.target
