# systemd/drukmix.service
[Unit]
Description=DrukMix Moonraker Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=dan

# Шлях до cfg (агент читає його при старті і при reload по mtime)
Environment=DRUKMIX_CONFIG=/home/dan/printer_data/config/drukmix.cfg

# Щоб Python не буферив stdout/stderr (краще видно в journal)
Environment=PYTHONUNBUFFERED=1

WorkingDirectory=/home/dan
ExecStart=/usr/bin/python3 /home/dan/drukmix/drukmix_agent.py

Restart=always
RestartSec=0.5

# Мінімальна “гігієна” без ломання serial доступу
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
