# systemd/drukmix.service
[Unit]
Description=DrukMix agent
After=network-online.target moonraker.service
Wants=network-online.target
Requires=moonraker.service

[Service]
Type=simple
User=drukos
WorkingDirectory=/home/drukos/drukmix
ExecStart=/home/drukos/drukmix/.venv/bin/python /home/drukos/drukmix/drukmix_agent.py --config /home/drukos/printer_data/config/drukmix.cfg
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
