# systemd/drukmix.service
[Unit]
Description=DrukMix driver
After=network-online.target moonraker.service
Wants=network-online.target
Requires=moonraker.service

[Service]
Type=simple
User=drukos
WorkingDirectory=/home/drukos/drukmix
Environment=DRUKMIX_CONFIG=/home/drukos/printer_data/config/drukmix_driver.cfg
ExecStart=/home/drukos/drukmix/.venv/bin/python /home/drukos/drukmix/drukmix_driver.py
Restart=always
RestartSec=2
Nice=5

[Install]
WantedBy=multi-user.target
